# target/systemd/user/bluez-volume.service
[Unit]
Description=bluez pulse volume setter
After=pulseaudio.socket

[Service]
Type=simple
ExecStart=%h/tools/bluez-volume-ng.sh
Restart=always

[Install]
WantedBy=default.target
