# target/systemd/user/gmediarender.service
[Unit]
Description=GMediaRender UPnP/DLNA renderer
After=pulseaudio.socket

[Service]
Type=simple
# %H hostname
ExecStart=/usr/bin/gmediarender -f %H --gstout-audiosink alsasink

[Install]
WantedBy=default.target
