# target/systemd/user/mini-agent.service
[Unit]
Description=bluez mini-agent

[Service]
Type=simple
ExecStart=%h/tools/mini-agent.py
Restart=always

[Install]
WantedBy=default.target
